# lab-2-build/mxnet/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:jupyter]
command=/usr/local/bin/jupyter-notebook --no-browser